# aws_deployment/crontab
@reboot cd node_mosquitto/ && npm start > /tmp/output.log
